# packaging/systemd/jw-authd@.service
[Unit]
Description=JW Agent one-shot PAM authentication worker
Requires=jw-authd.socket
CollectMode=inactive-or-failed

[Service]
Type=simple
User=root
Group=root
EnvironmentFile=-/etc/default/jw-agent
ExecStart=/usr/lib/jw-agent/jw-authd
StandardInput=socket
StandardOutput=socket
StandardError=journal
TimeoutStartSec=12s
UMask=0077
NoNewPrivileges=yes
PrivateDevices=yes
PrivateTmp=yes
ProtectClock=yes
ProtectControlGroups=yes
ProtectHome=read-only
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectSystem=strict
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native
