# packaging/systemd/jw-agentd.service
[Unit]
Description=JW Agent unprivileged management daemon
After=local-fs.target jw-authd.socket jw-opsd.service
Wants=jw-authd.socket jw-opsd.service

[Service]
Type=simple
User=jw-agent
Group=jw-agent
SupplementaryGroups=jw-agent-proxy
EnvironmentFile=-/etc/default/jw-agent
ExecStart=/usr/lib/jw-agent/jw-agentd serve
Restart=on-failure
RestartSec=2s
UMask=0077
NoNewPrivileges=yes
PrivateDevices=yes
PrivateTmp=yes
ProtectClock=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/jw-agent/agentd /run/jw-agent /run/jw-agent-proxy
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
