# packaging/systemd/jw-opsd.service
[Unit]
Description=JW Agent networkless privileged operation helper
After=local-fs.target

[Service]
Type=simple
User=root
Group=root
EnvironmentFile=-/etc/default/jw-agent
ExecStart=/usr/lib/jw-agent/jw-opsd
Restart=on-failure
RestartSec=2s
UMask=0007
NoNewPrivileges=yes
PrivateDevices=yes
PrivateNetwork=yes
PrivateTmp=yes
ProtectClock=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectSystem=strict
ReadWritePaths=/run/jw-agent /var/lib/jw-agent/opsd /etc/nginx/sites-enabled /etc/nginx/sites-available /etc/php/8.3/fpm
BindPaths=/dev/null:/var/log/php8.3-fpm.log
TemporaryFileSystem=/var/log/nginx:rw
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
IPAddressDeny=any
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
