# packaging/systemd/jw-opsd.service
[Unit]
Description=JW Agent networkless privileged operation helper
After=local-fs.target

[Service]
Type=simple
User=root
Group=jw-agent
EnvironmentFile=-/etc/default/jw-agent
ExecStart=/usr/lib/jw-agent/jw-opsd
Restart=on-failure
RestartSec=2s
UMask=0007
NoNewPrivileges=yes
PrivateDevices=yes
PrivateTmp=yes
ProtectClock=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectSystem=strict
ReadWritePaths=/run/jw-agent
RestrictAddressFamilies=AF_UNIX
IPAddressDeny=any
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
